# ansible/ansible_install_k8s_v3/roles/certs-renew/templates/k8s-certs-renew.service.j2
[Unit]
Description=Renew K8S control plane certificates

[Service]
Type=oneshot
ExecStart=/etc/kubernetes/k8s-certs-renew.sh