# ansible/ansible_install_k8s_v2/roles/certs-renew/templates/k8s-certs-renew.timer.j2
[Unit]
Description=Timer to renew K8S control plane certificates

[Timer]
OnCalendar=*-10-15 03:{{ groups['kube_control_plane'].index(inventory_hostname) }}0:00

[Install]
WantedBy=multi-user.target